# dockerize/supervisord.conf
[supervisord]
nodaemon=true

[program:update-blog-app]
command = python /var/www/blog/myblog-update/main.py

[program:nginx-app]
command = /usr/sbin/nginx -g "daemon off;"
